# crontab
# Run every minute
*/1 * * * * cf --version

# Run every 2 seconds
*/2 * * * * * * echo "hello world, every 2 seconds"